# task/5808/unuse/bin/network_watch
#!/bin/sh /etc/rc.common

{
	local lanif
	local linkmode
	local prev_linkmode
	local wlan_mode
	local wlan_up
	local wlan_led_state
	local lte_led_state
	local lte_link
	local str
	config_load network
	config_get lanif lan ifname
	
	wlan_mode=""
	wlan_up="no"
	wlan_led_state="off"
	lte_led_state=""
	config_load wireless
	config_get wlan_mode wifi0 mode

	lanif=eth0.1
	#echo $lanif > /dev/console
	[ -z "$lanif" ] && exit 0
                                              
	prev_linkmode=yes                         
	while true                                
	do                                        
		linkmode=$(ethtool $lanif)        
		linkmode=${linkmode#*Link detected: }
		linkmode=${linkmode%%\n}
		while [ "$linkmode" != "yes" -a "$linkmode" != "no" ]
		do                 
			sleep 1    
			linkmode=$(ethtool $lanif)
			linkmode=${linkmode#*Link detected: }
			linkmode=${linkmode%%\n}
		done                      
		#echo $linkmode > /dev/console

		[ -n "$linkmode" ] && [ "$linkmode" != "$prev_linkmode" ] && {
			prev_linkmode=$linkmode                      

			[ "$linkmode" = "no" ] && {
				echo link down
			}                       
			[ "$linkmode" = "yes" ] && {
				echo link up
			}                                         
			#/etc/init.d/network reload 
		}

		str=""
		str=`ifconfig | grep "ra0"`
		if [ -n "$str" ]; then
			wlan_up="yes"
		else
			wlan_up="no"
		fi
		if [ -n "$wlan_mode" ] && [ "$wlan_mode" = "ap" ] && [ "$wlan_up" = "yes" ]; then
			if [ "$wlan_led_state" = "off" ]; then
				echo "timer" > /sys/class/leds/wlan-run/trigger
				echo "500" > /sys/class/leds/wlan-run/delay_on
				echo "500" > /sys/class/leds/wlan-run/delay_off
				wlan_led_state="on"
			fi
		elif [ "$wlan_led_state" != "off" ]; then
			echo "none" > /sys/class/leds/wlan-run/trigger
			echo "255" > /sys/class/leds/wlan-run/brightness
			wlan_led_state="off"
		fi
		
		lte_link=$(ifconfig 3g-wan2 2> /dev/null)
		if [ -n "$lte_link" -a "$lte_led_state" != "on" ]; then
			echo "0" > /sys/class/leds/lte-run/brightness
			lte_led_state="on"
		elif [ -z "$lte_link" -a "$lte_led_state" != "off" ]; then
			echo "255" > /sys/class/leds/lte-run/brightness
			lte_led_state="off"
		fi
		sleep 1                                           
	done        
}
