# task/5808/bin/restore_default
#!/bin/sh /etc/rc.common

echo "restore default" > /dev/console
[ -f /etc/config/wireless ] && {
	local wifimode
	config_load wireless
	config_get wifimode wifi0 mode
	[ -n "$wifimode" ] && [ "$wifimode" = "sta" ] && {
		echo > /etc/config/network
		cat > /etc/config/network << EOF
config interface 'loopback'
        option ifname 'lo'
        option proto 'static'
        option ipaddr '127.0.0.1'
        option netmask '255.0.0.0'

config globals 'globals'
        option ula_prefix 'fdca:5f92:f802::/48'

config interface 'lan'
        option ifname 'eth0.1'
        option force_link '1'
        option macaddr '00:11:22:33:44:55'
        option type 'bridge'
        option proto 'static'
        option ipaddr '192.168.11.1'
        option netmask '255.255.255.0'
        option ip6assign '60'
        option disabled '0'

config interface 'wan'
        option ifname 'ra0'
        option force_link '1'
        option macaddr '00:11:22:33:44:56'
        option proto 'dhcp'
EOF
	}
	[ -n "$wifimode" ] && [ "$wifimode" = "ap" ] && {
		echo > /etc/config/network
		cat > /etc/config/network << EOF
config interface 'loopback'
        option ifname 'lo'
        option proto 'static'
        option ipaddr '127.0.0.1'
        option netmask '255.0.0.0'

config globals 'globals'
        option ula_prefix 'fdca:5f92:f802::/48'

config interface 'eth0'
        option ifname 'eth0'
        option macaddr '00:11:22:33:44:55'

config interface 'lan'
        option ifname 'eth0.1'
        option force_link '1'
        option macaddr '00:11:22:33:44:55'
        option type 'bridge'
        option proto 'static'
        option ipaddr '192.168.11.1'
        option netmask '255.255.255.0'
        option ip6assign '60'
        option disabled '0'

config interface 'wan'
        option ifname 'eth0.2'
        option force_link '1'
        option macaddr '00:11:22:33:44:56'
        option proto 'dhcp'

config interface 'wan6'
        option ifname 'eth0.2'
        option proto 'dhcpv6'
EOF
	}
}
rm -f /etc/config/profile_network
rm -f /etc/config/network_tmp
/etc/init.d/extension start
/etc/init.d/network restart
uci set lucid.http.address=80
uci add_list lucid.http.address=8345
uci add_list lucid.http.address=8848
uci set lucid.https.address=443
uci commit lucid
uci set system.telnet.port=23
uci commit system
uci set dropbear.main.Port=22
uci commit dropbear
lua /usr/lib/lua/luci/scripts/reset_access_control.lua
/etc/init.d/lucid restart
/etc/init.d/telnet restart
/etc/init.d/dropbear restart
(echo "admin";sleep 1;echo "admin")|passwd admin > /dev/null

echo "none" > /sys/class/leds/uc100-run/trigger
echo "0" > /sys/class/leds/uc100-run/brightness
